# face/supervisord.conf
[supervisord]
nodaemon=true

[program:frontend]
command=/usr/local/bin/python3 /python/app/entrypoint.py
