# face/supervisord.conf
[supervisord]
nodaemon=true

[program:face]
command=/usr/local/bin/python3 /python/app/run.py
